# sv/dhcpd3/run
#!/bin/sh
set -e
exec 2>&1

echo "**** Starting: dhcpd3"
INTERFACES="$(cat conf/INTERFACES)"
exec dhcpd3 -f -d $INTERFACES
