# sv/apache2/run
#!/bin/sh

echo "**** Starting: apache2"
exec apache2 -DNO_DETACH 2>&1
