# sv/lighttpd/run
#!/bin/sh
set -e

exec lighttpd -f /etc/lighttpd/lighttpd.conf -D 2>&1
