# sv/smartd/run
#!/bin/sh

echo '**** Starting: approx'
exec smartd -d --interval=3600 2>&1
