# sv/smartd/run
#!/bin/sh

exec smartd -d --interval=3600
