# sv/mdadm/run
#!/bin/sh
set -e

exec mdadm --monitor --scan --syslog 2>&1
