# sv/mdadm/run
#!/bin/sh
exec mdadm --monitor --scan --syslog 2>&1
