# sv/syslog-ng/run
#!/bin/sh
set -e

exec syslog-ng -F 2>&1
