# sv/cupsd/run
#!/bin/sh
set -e

chown root:lpadmin /usr/share/cups/model
chmod 3775 /usr/share/cups/model
exec cupsd -f
