# sv/git-daemon/run
#!/bin/sh
set -e

echo '**** Starting: git-daemon'
exec git daemon --verbose --export-all --base-path=/var/cache /var/cache/git 2>&1
