# sv/ssh/run
#!/bin/sh
set -e

# sshd wants to be started with absolute pathname
SSHD=$(which sshd) || sv down .

PRIVSEP_DIR=/var/run/sshd
mkdir -p "$PRIVSEP_DIR"

exec "$SSHD" -D 2>&1
