# sv/postfix/run
#!/bin/sh
set -e

postfix check
exec /usr/lib/postfix/master 1>&2
