# sv/pywebdav/run
#!/bin/sh
set -e
exec 2>&1

NAME="$(basename "$(pwd)")"
echo "**** Starting: $NAME"

exec chpst -u "www-data" davserver --directory /var/www --noauth --nomime --host localhost --port 8008
