# sv/approx/run
#!/bin/sh

echo 'approx starting'
exec approx -f 2>&1
