# sv/approx/run
#!/bin/sh

echo '**** Starting: approx'
exec approx -f 2>&1
