# sv/ntpd/run
#!/bin/sh
set -e

exec ntpd -u ntp -g -n 2>&1
